# local/etc/systemd/syslog-krb5.timer
[Unit]
Description=Run 'syslog-krb5.service' four hourly and on boot

[Timer]
OnUnitActiveSec=4h

[Install]
WantedBy=timers.target
